# systemd/oysters-dump.timer
[Unit]
Description=Run Oysters dump daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
